# spec/fixtures/dev.elasticsearch.plist
<?xml version="1.0" encoding="UTF-8"?>

<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>dev.elasticsearch</string>

    <key>Boxen</key>
    <dict>
      <key>Description</key>
      <string>Full-text indexing and searching.</string>

      <key>Port</key>
      <integer>19200</integer>
    </dict>

    <key>ProgramArguments</key>
    <array>
      <string>/test/boxen/homebrew/bin/elasticsearch</string>
      <string>-f</string>
      <string>-D es.config=/test/boxen/config/elasticsearch/elasticsearch.yml</string>
    </array>

    <key>RunAtLoad</key>
    <true />

    <key>KeepAlive</key>
    <true />

    <key>UserName</key>
    <string>testuser</string>

    <key>WorkingDirectory</key>
    <string>/test/boxen/data/elasticsearch</string>

    <key>StandardErrorPath</key>
    <string>/test/boxen/log/elasticsearch/console.log</string>

    <key>StandardOutPath</key>
    <string>/test/boxen/log/elasticsearch/console.log</string>
  </dict>
</plist>
